# templates/systemd/netops-ping.service
[Unit]
Description=NetOps Ping Engine
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/netops-box
ExecStart=/opt/netops-box/.venv/bin/python /opt/netops-box/modules/ping/ping_engine.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
